# files/etc/cron.d/rsync
# Run backup script
*/5 * * * * root /usr/local/bin/rsync.sh > /dev/null
